# client/homepost-client.service
[Unit]
Description=Home Monitoring Client
After=network.target

[Service]
Type=simple
User=%USER%
WorkingDirectory=%INSTALL_DIR%
ExecStart=/usr/bin/env node client.js
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=homepost-client

[Install]
WantedBy=multi-user.target
